# scripts/dewey/hermes-voice-bridge.supervisor.conf
[program:hermes-voice-bridge]
directory=/root/Desktop/repos/hermes-voice-ios
command=/root/Desktop/repos/hermes-voice-ios/scripts/dewey/run-bridge.sh
user=root
autostart=true
autorestart=unexpected
startsecs=5
startretries=3
stopsignal=TERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true
environment=HOME="/root",PATH="/root/.hermes/node/bin:/usr/local/bin:/usr/bin:/bin"
stdout_logfile=/var/log/hermes-voice-bridge.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/hermes-voice-bridge.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
